# server/templates/weatherstation.cron.template
* * * * * {{ username }} {{ env_path }}/python {{ repo_path }}/server/jobs/fetch.py
