# server/templates/weatherstation.cron.template
* * * * * {{ username }} cd {{repo_path}}/server && {{ env_path }}/python jobs/fetch.py
